# meta-refkit-core/recipes-security/nftables-settings-default/files/firewall.path
[Unit]
Description=Track changes to firewall configuration

[Path]
ReloadOnTrigger=true
PathChanged=/run/firewall/firewall.ruleset

[Install]
WantedBy=network.target
